# systemd/cachyos-update.service
[Unit]
Description=Actualiza el sistema CachyOS (pacman + yay) y notifica
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
ExecStart=@SCRIPTS_DIR@/update-system.sh
