# systemd/omarchy-check.service
[Unit]
Description=Comprueba nuevas versiones de omarchy-on-cachyos
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
ExecStart="@SCRIPTS_DIR@/check-omarchy-update.sh"
